# etc/systemd/system/storagefancontrol.service
[Unit]
Description=storagefancontrol

[Service]
Type=simple
ExecStart=/usr/sbin/storagefancontrol

[Install]
WantedBy=multi-user.target

